# build/kurohai-tunnels.service
[Unit]
Description=Kurohai's Tunnel Service
After=network-online.target


[Service]
Type=forking
Restart=always
Environment="RSA_KEY=/etc/kurohai/ssh/id_rsa"
Environment="CONFIG_DIR=/etc/kurohai/tunnels.d"
Environment="LOG_DIR=/var/log/autossh"
Environment="PID_DIR=/var/run/autossh"
ExecStart=/etc/kurohai/start-tunnels.sh
ExecStop=/bin/kill -s TERM $MAINPID


[Install]
WantedBy=multi-user.target
