# build/kurohai-hdmi-controller-app.service
[Unit]
Description=Kurohai's HDMI Controller Service
After=network-online.target


[Service]
Type=forking
Environment="VIRTUAL_ENV=/home/pi/.virtualenvs/hdmi-control-ui"
Environment="WORKON_HOME=/home/pi/.virtualenvs"
; Environment="LOG_DIR=/var/log/autossh"
; Environment="PID_DIR=/var/run/autossh"
ExecStart=/etc/kurohai/start-app.sh
; ExecStop=pkill -9 -f ".*uwsgi.*--ini.*uwsgi-hdmi-control-ui-devel.ini"
ExecStop=/bin/kill -s TERM $MAINPID


[Install]
WantedBy=multi-user.target
